# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
childlogdir = /var/log/supervisor

[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock         ; use a unix:// URL  for a unix socket
loglevel=debug

[program:initconfig]
command=/usr/local/sbin/fix_certs.sh
autorestart=false
startretries=0
startsecs=0
redirect_stderr=false
redirect_stderr=false

[program:redirector-xrootd]
command=xrootd -c /etc/xrootd/cache-redirector.cfg -k fifo -n cache-redirector -k 10 -s /var/run/xrootd/cache-redirector.pid -l /var/log/xrootd/xrootd.log
user=xrootd
autorestart=true

[program:redirector-cmsd]
command=/usr/bin/cmsd -c /etc/xrootd/cache-redirector.cfg -p 9004 -n cache-redirector -k 7 -s /var/run/xrootd/cache-redirector.pid -l /var/log/xrootd/cmsd.log
user=xrootd
autorestart=true

[program:crond]
command=/usr/sbin/crond -n
autorestart=true
